# docker/etc/supervisord.d/processes.ini
[program:gunicorn]
command=/usr/bin/gunicorn -c /opt/netbox/gunicorn_config.py netbox.wsgi
directory=/opt/netbox/netbox/
user=nginx
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
